# deploy/stofgard-worker.conf
[program:stofgard-worker]
process_name=%(program_name)s_%(process_num)02d
command=docker exec stofgard-laravel.test-1 php artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/stofgard-worker.log
stopwaitsecs=3600
